# second-brain.service
[Unit]
Description=Second Brain Telegram Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/second-brain-bot
ExecStart=/home/ubuntu/.venv/bin/python main.py
Restart=always
RestartSec=10
EnvironmentFile=/home/ubuntu/second-brain-bot/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
