# second-brain.service
[Unit]
Description=Second Brain Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/second-brain-bot
ExecStart=/root/second-brain-bot/.venv/bin/python main.py
Restart=always
RestartSec=10
EnvironmentFile=/root/second-brain-bot/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
